# supervisord.conf
[supervisord]
nodaemon=true

[program:timemachine]
command=/start_netatalk.sh