# supervisord.conf
[supervisord]
nodaemon=true

[program:timemachine]
command=/entrypoint.sh